# firewall-auth.conf
description "Firewall authentication script"

start on startup
script
  /root/firewall-auth.sh
end script
